# deploy/systemd/tender-monitor-scheduler.service
[Unit]
Description=Tender Monitor — connector scheduler (APScheduler, runs ingest jobs per source)
Documentation=https://github.com/your-org/tender-monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tender
Group=tender
WorkingDirectory=/opt/tender-monitor/app
ExecStart=/opt/tender-monitor/venv/bin/python -m tender_monitor.cli run-scheduler

Restart=on-failure
RestartSec=10
# Give in-flight ingests time to finish so the DB doesn't end up in a half state.
# 60s is generous — covers a slow national_bank detail-page fetch with margin.
TimeoutStopSec=60
KillSignal=SIGINT

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tender-monitor-scheduler

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
