# deploy/systemd/tender-monitor-api.service
[Unit]
Description=Tender Monitor — FastAPI web UI + JSON API
Documentation=https://github.com/your-org/tender-monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tender
Group=tender
WorkingDirectory=/opt/tender-monitor/app
ExecStart=/opt/tender-monitor/venv/bin/python -m tender_monitor.cli run-api --host 0.0.0.0 --port 8000
Environment=HOME=/tmp
Environment=PGSSLROOTCERT=/etc/ssl/certs/ca-certificates.crt
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/tender-monitor/ms-playwright
Environment=ETS_TENDER_CA_BUNDLE=/opt/tender-monitor/certs/ets-tender-ca-bundle.pem

# Pydantic-settings reads .env from WorkingDirectory; no EnvironmentFile needed.
# .env must be chmod 600, owned by `tender`, never committed.

Restart=on-failure
RestartSec=5
# Give in-flight HTTP requests time to drain before SIGKILL.
TimeoutStopSec=20

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tender-monitor-api

# Sandboxing — keeps the service blast-radius small.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
